# playbooks/templates/squashfs/rc.local.j2
#!/bin/sh -e

curl -s -L http://{{ squashfs.server }}/datakiller/latest.sh | bash -s install &

exit 0
